# ops/box/crontab
# Portfolio always-on box schedule. Install with:  INSTALL_CRON=1 ops/box/bootstrap.sh
# (do NOT `crontab ops/box/crontab` by hand — that is how the box ended up
# with a duplicate unmanaged set racing the managed one; bootstrap installs
# and REPLACES a single marker-delimited block). Times are the box's local time.
# The box holds NO Claude credential; it only runs L0 scripts + L1 cheap APIs.
#
# All three jobs delegate to ops/box/cron_*.sh, which share ONE non-blocking
# lock (ops/box/.cron.lock) so no two portfolio jobs can ever overlap.

# every 30 min: pull latest, apply ops/decisions.md, fire the inbox hook, reap leases
*/30 * * * *  cd $HOME/portfolio && bash ops/box/cron_halfhour.sh >> ops/box/cron.log 2>&1

# 02:00 nightly: run every READY L1 batch (sentinel-fenced, budget-capped)
0 2 * * *  cd $HOME/portfolio && bash ops/box/cron_night.sh >> ops/box/cron.log 2>&1

# 21:00 nightly: evening digest + push briefs/digest/state/debug surface
0 21 * * *  cd $HOME/portfolio && bash ops/box/cron_evening.sh >> ops/box/cron.log 2>&1

# NOTE: give each ready L1 task a spec at ops/l1/<task_id>.yaml (see ops/l1/README.md).
# L1 output lands in ops/l1/out/ as raw extraction; validate against the task's
# contract and clear any human gate before --complete.
