# research-portfolio/ops/box/crontab
# Portfolio always-on box schedule. Install with:  crontab ops/box/crontab
# Edit REPO below to your clone path. Times are the box's local time.
# The box holds NO Claude credential; it only runs L0 scripts + L1 cheap APIs.

# --- shared prologue: cd into the repo, load L1 keys, use the venv python ---
# (cron has a minimal env, so we set PATH and source the box .env explicitly)

# every 30 min: expire stale leases so abandoned tasks re-queue
*/30 * * * *  cd $HOME/portfolio/research-portfolio && .venv/bin/python ops/runner/runner.py --reap >> ops/box/cron.log 2>&1

# 02:00 nightly: run the overnight L1 batch layer (sentinel-fenced, budget-capped)
0 2 * * *  cd $HOME/portfolio/research-portfolio && set -a && . ops/box/.env && set +a && .venv/bin/python ops/runner/dispatch.py --smoke >> ops/box/cron.log 2>&1

# 21:00 nightly: write the evening digest, then push briefs + digest to GitHub
0 21 * * *  cd $HOME/portfolio/research-portfolio && .venv/bin/python ops/runner/runner.py --digest >> ops/box/cron.log 2>&1 && git add ops/briefs ops/digest && git commit -q -m "box: nightly digest $(date +\%F)" && git push -q origin HEAD:main >> ops/box/cron.log 2>&1

# NOTE: models.py is wired for live dispatch. The 02:00 line runs --smoke as a
# nightly health check; run real batches with:
#   dispatch.py --worker <w> --items <batch.jsonl> --live --cost <N> --out <answers.json>
# Once a per-task L1 driver exists (iterate the plan's ready L1 set -> build items
# -> dispatch -> --complete on sentinel pass), swap the 02:00 line for it.
